# deploy/leap.service
[Unit]
Description=Leap Trading API Server
After=network.target xtqmt.service
Wants=xtqmt.service

[Service]
Type=simple
User=leap
Group=leap
WorkingDirectory=/opt/leap/current
Environment=PYTHONUTF8=1
Environment=WINEPREFIX=/opt/leap/wine
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/bash -c 'while ! pgrep -f "XtMiniQmt.exe" >/dev/null 2>&1; do sleep 2; done'
ExecStart=/usr/bin/wine python -m uvicorn leap.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target